# config/systemd/canyougrab-worker@.service
[Unit]
Description=CanYouGrab Domain Worker %i
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/canyougrab/api
EnvironmentFile=/opt/canyougrab/database.env
EnvironmentFile=/opt/canyougrab/valkey.env
Environment=BATCH_CONCURRENCY=10
ExecStart=/opt/canyougrab/venv/bin/python3 /opt/canyougrab/api/worker.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=canyougrab-worker-%i

[Install]
WantedBy=multi-user.target
