# config/systemd/rq-metrics.service
[Unit]
Description=CanYouGrab RQ Metrics Exporter (multi-environment)
After=network.target

[Service]
Type=simple
User=root
EnvironmentFile=/opt/canyougrab-admin/valkey.env
Environment=RQ_METRICS_PORT=9122
Environment=RQ_QUEUE_ENVIRONMENTS=dev:queue:jobs:dev,prod:queue:jobs:prod
Environment=RQ_HOST_ENVIRONMENTS=dev:dev,prod:prod
ExecStart=/opt/canyougrab-admin/venv/bin/python /opt/canyougrab-repo/scripts/rq_metrics_exporter.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
